# nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=3